# agents/artifacts/TASK-018-supervisor-specific-launcher-packaging-and-stop-contract/task_scheduler_definition.xml
<?xml version="1.0" encoding="UTF-16"?>

<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>

    <Date>2026-04-08T08:10:55</Date>

    <Author>DESKTOP-PJC\PJC</Author>

    <URI>\SOOP-LiveVOD-Task018</URI>

  </RegistrationInfo>

  <Principals>

    <Principal id="Author">

      <UserId>S-1-5-21-3197170921-3206713947-2637509415-1001</UserId>

      <LogonType>InteractiveToken</LogonType>

    </Principal>

  </Principals>

  <Settings>

    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>

    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>

    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>

    <IdleSettings>

      <Duration>PT10M</Duration>

      <WaitTimeout>PT1H</WaitTimeout>

      <StopOnIdleEnd>true</StopOnIdleEnd>

      <RestartOnIdle>false</RestartOnIdle>

    </IdleSettings>

  </Settings>

  <Triggers>

    <TimeTrigger>

      <StartBoundary>2099-01-01T00:00:00</StartBoundary>

    </TimeTrigger>

  </Triggers>

  <Actions Context="Author">

    <Exec>

      <Command>powershell.exe</Command>

      <Arguments>-NoProfile -ExecutionPolicy Bypass -File C:\python\live_vod_scheduler_task_action.ps1 -RuntimeDir C:\python\data\t018rt -DbPath C:\python\data\t018.db -Port 49728 -HealthTimeout 15</Arguments>

    </Exec>

  </Actions>

</Task>